# aria2/.config/aria2/aria2.service
[Install]
WantedBy=default.target

[Service]
ExecStart=aria2c
Restart=on-failure

[Unit]
After=network.target
Description=aria2
